# apache/supervisor.cosmohome.conf
[program:certbot]
command=sh -c "certbot run -n --apache -d cosmologyathome.org,www.cosmologyathome.org"
autorestart=false
umask=0002
